# www-misc/profile-sync-daemon/files/psd
#!/sbin/runscript

depend() {
	need localmount
	use cron
}

start() {
	ebegin Starting Profile-Sync-Daemon
	/usr/bin/profile-sync-daemon sync
	eend $?
}

stop() {
	ebegin Stopping Profile-Sync-Daemon
	/usr/bin/profile-sync-daemon unsync
	eend $?
}
